# share/system_files/etc/profile.d/renkei-vpe.csh
setenv RVPE_LOCATION /usr/rpop/renkei-vpe
setenv PATH $RVPE_LOCATION/bin:$RVPE_LOCATION/sbin:$PATH
