# share/system_files/etc/profile.d/renkei-vpe.sh
export RVPE_LOCATION=/usr/rpop/renkei-vpe
export PATH=$RVPE_LOCATION/bin:$RVPE_LOCATION/sbin:$PATH
export MANPATH=$RVPE_LOCATION/share/man:$MANPATH
export RVPE_XMLRPC=http://localhost:3111/
